# fedora/rootfs/usr/lib/systemd/system/octoprint-firstboot.service
[Unit]
Description=Initalize octoprint container

[Service]
Type=oneshot
ExecStart=/usr/bin/octoprint-bootstrap

[Install]
WantedBy=multi-user.target